# VM/supervisor.conf
[supervisord]
nodaemon=true

[program:php-fpm]
command=php-fpm8.0 -R -F -c /etc/php/8.0/fpm/php-fpm.conf
stdout_logfile=/var/log/php-fpm/stdout.log
stdout_logfile_maxbytes=0
stderr_logfile=/var/log/php-fpm/stderr.log
stderr_logfile_maxbytes=0
autostart=true
autorestart=unexpected

[program:node]
directory=/var/www/html/main/frontend
autostart=true
command=npm run start